# systemd/palworld-restarter.timer
[Unit]
  Description=Palworld Server Restart Timer

[Timer]
  OnCalendar=*-*-* 06,18:00:00
  Persistent=true

[Install]
  WantedBy=timers.target